# scripts/upstart/pingd.conf
description "ping daemon"
author "Chris Fallin <cfallin@c1f.net>"

start on net-service-up IFACE!=lo
stop on runlevel [!2345]

exec /usr/local/bin/pingd.sh
